# deploy/infonigma.socket
[Unit]
Description=infonigma socket

[Socket]
ListenStream=/run/infonigma/socket

[Install]
WantedBy=sockets.target
